# deploy/transcriber.service
[Unit]
Description=Transcriber Web Server
After=network.target

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/opt/transcriber
Environment="PATH=/home/ec2-user/.local/bin:/usr/local/bin:/usr/bin"

# 清掉上次 crash 殘留的 multiprocessing 子進程。在 master 啟動「前」跑一次，
# 而不是每個 worker startup hook 跑一次。
# 原本放在 src/main.py 的 startup，但 uvicorn --workers 2 之下每個 worker 都
# 會跑 pkill，把 sibling worker 的 multiprocessing.resource_tracker 一起殺掉
# → workers crash loop。把這步移到 systemd ExecStartPre，整個 service 只跑一次。
#
# 走獨立 script 而非 inline /bin/sh -c '...'，原因見 pre-start-cleanup.sh 內註解
# （簡言之：sh -c 的參數會出現在 sh 進程 cmdline 內，pkill -f 會把 sh 自己當 match 殺掉）。
ExecStartPre=/opt/transcriber/deploy/pre-start-cleanup.sh

# 綁 127.0.0.1：nginx 走 loopback 連後端，uvicorn 不需對外。綁 0.0.0.0 會讓
# 外網經 SG 開放的 8000 直接打進後端、繞過 nginx 的 rate limit / CSP / docs 關閉。
ExecStart=/home/ec2-user/.local/bin/uvicorn src.main:app --host 127.0.0.1 --port 8000 --workers 2
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
