# deploy/transcriber.service
[Unit]
Description=Transcriber Web Server
After=network.target

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/opt/transcriber
Environment="PATH=/home/ec2-user/.local/bin:/usr/local/bin:/usr/bin"

# 清掉上次 crash 殘留的 multiprocessing 子進程。在 master 啟動「前」跑一次，
# 而不是每個 worker startup hook 跑一次。
# 原本放在 src/main.py 的 startup，但 uvicorn --workers 2 之下每個 worker 都
# 會跑 pkill，把 sibling worker 的 multiprocessing.resource_tracker 一起殺掉
# → workers crash loop。把這步移到 systemd ExecStartPre，整個 service 只跑一次。
ExecStartPre=/usr/bin/sh -c 'pkill -9 -f multiprocessing.spawn || true; pkill -9 -f multiprocessing.resource_tracker || true; exit 0'

ExecStart=/home/ec2-user/.local/bin/uvicorn src.main:app --host 0.0.0.0 --port 8000 --workers 2
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
